# files/etc/init.d/crowsnest
#!/bin/sh /etc/rc.common

START=99
STOP=15

EXTRA_COMMANDS="status"
EXTRA_HELP="        status          Check service status"

start() {
    echo "Starting TollGate crowsnest" > /tmp/crowsnest.log 2>&1
    tollgate-crowsnest >> /tmp/crowsnest.log 2>&1 &
    echo "PID: $!" >> /tmp/crowsnest.log 2>&1
}

stop() {
    echo "Stopping TollGate crowsnest"
    killall tollgate-crowsnest
}

status() {
    if pgrep -f "tollgate-crowsnest" >/dev/null; then
        echo "TollGate crowsnest is running"
        return 0
    else
        echo "TollGate crowsnest is not running"
        return 1
    fi
}
