# files/etc/init.d/merchant
#!/bin/sh /etc/rc.common

START=99
STOP=15

EXTRA_COMMANDS="status"
EXTRA_HELP="        status          Check service status"

start() {
    echo "Starting TollGate merchant" > /tmp/merchant.log 2>&1
    tollgate-merchant >> /tmp/merchant.log 2>&1 &
    echo "PID: $!" >> /tmp/merchant.log 2>&1
}

stop() {
    echo "Stopping TollGate merchant"
    killall tollgate-merchant
}

status() {
    if pgrep -f "tollgate-merchant" >/dev/null; then
        echo "TollGate merchant is running"
        return 0
    else
        echo "TollGate merchant is not running"
        return 1
    fi
}
