# files/etc/init.d/relay
#!/bin/sh /etc/rc.common

START=99
STOP=15

EXTRA_COMMANDS="status"
EXTRA_HELP="        status          Check service status"

start() {
    echo "Starting TollGate relay" > /tmp/relay.log 2>&1
    tollgate-relay >> /tmp/relay.log 2>&1 &
    echo "PID: $!" >> /tmp/relay.log 2>&1
}

stop() {
    echo "Stopping TollGate relay"
    killall tollgate-relay
}

status() {
    if pgrep -f "tollgate-relay" >/dev/null; then
        echo "TollGate relay is running"
        return 0
    else
        echo "TollGate relay is not running"
        return 1
    fi
}
